# mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/additional-prompts.conf
[Service]
ExecStart=
ExecStart=systemd-firstboot --prompt-locale --prompt-keymap --prompt-root-password --mute-console=yes
